# contrib/init/ellocashd.conf
description "Ellocash Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env ELLOCASHD_BIN="/usr/bin/ellocashd"
env ELLOCASHD_USER="ellocash"
env ELLOCASHD_GROUP="ellocash"
env ELLOCASHD_PIDDIR="/var/run/ellocashd"
# upstart can't handle variables constructed with other variables
env ELLOCASHD_PIDFILE="/var/run/ellocashd/ellocashd.pid"
env ELLOCASHD_CONFIGFILE="/etc/ellocash/ellocash.conf"
env ELLOCASHD_DATADIR="/var/lib/ellocashd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # ellocashd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$ELLOCASHD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run ellocashd."
        echo "The setting must appear in $ELLOCASHD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Ellocash Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$ELLOCASHD_PIDDIR"
    chmod 0755 "$ELLOCASHD_PIDDIR"
    chown $ELLOCASHD_USER:$ELLOCASHD_GROUP "$ELLOCASHD_PIDDIR"
    chown $ELLOCASHD_USER:$ELLOCASHD_GROUP "$ELLOCASHD_CONFIGFILE"
    chmod 0660 "$ELLOCASHD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$ELLOCASHD_PIDFILE" \
    --chuid $ELLOCASHD_USER:$ELLOCASHD_GROUP \
    --exec "$ELLOCASHD_BIN" \
    -- \
    -pid="$ELLOCASHD_PIDFILE" \
    -conf="$ELLOCASHD_CONFIGFILE" \
    -datadir="$ELLOCASHD_DATADIR" \
    -disablewallet \
    -daemon

